# systemd/scheduler-seriea.service
[Unit]
Description=Vertex Football Scheduler — Serie A
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vertex
WorkingDirectory=/opt/vertex-football-scraper
ExecStart=/opt/vertex-football-scraper/.venv/bin/python scheduler.py --league SERIEA
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=scheduler-seriea

[Install]
WantedBy=multi-user.target
